# Linux-config/etc/rc.local
#!/bin/sh

/usr/local/bin/shutdowncheck.py &

exit 0
